# pipelines/services/news/deploy/tinboker-news.service
[Unit]
Description=TinBoker news ingest (RSS feeds -> shared knowledge wiki)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/tinboker-agents/services/news
ExecStart=/root/tinboker-agents/services/news/scripts/run_news.sh
Environment=GOOGLE_APPLICATION_CREDENTIALS=/root/tinboker-agents/services/podcast/gcp-service-account.json
# A single run is bounded; cap it so a hung fetch cannot wedge the timer.
TimeoutStartSec=1800
