# backend/deploy/systemd/tinboker-staging.service
[Unit]
Description=TinBoker Staging Backend
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/app/backend
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.staging.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.staging.yml down
ExecReload=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.staging.yml up -d --force-recreate

[Install]
WantedBy=multi-user.target
